# bin/swsdvolcal
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swsdvolcal
# Required-Start:    $network $syslog
# Required-Stop:     $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: swsdvolcal on port 9004
# Description:       lists swsdvolcal 
### END INIT INFO

me=$(basename $0)

# Do NOT "set -e"
USER="www"
DIR="/usr/share/nginx/node/swsdvolcal"

if [ "$(id -un)" != "${USER}" ]
then
  su - "${USER}" -c "/etc/init.d/swsdvolcal $*"
  exit
fi

# PATH should only include /usr/* if it runs after the mountnfs.sh script
PATH="/sbin:/usr/sbin:/bin:/usr/bin:/usr/local/bin"
DESC="swsdvolcal"
NAME="swsdvolcal"
ROOT="${DIR}"
GROUP="www"
DAEMON="/usr/bin/node"
DAEMON_ARGS="${ROOT}/app.js"
PIDFILE="/var/opt/swsdvolcal/run/${NAME}.pid"
SCRIPTNAME="/etc/init.d/${NAME}"
export NODE_ENV="production"
LOGDIR="/var/log/node/${NAME}"

# Exit if the package is not installed
[ -x "${DAEMON}" ] || exit 1

# Read configuration variable file if it is present
[ -r "/etc/default/${NAME}" ] && . "/etc/default/${NAME}"

cd "${DIR}" || exit 2

# load nvm && forever
export NVM_DIR="/home/ec2-user/.nvm"
[ -s "$NVM_DIR/nvm.sh" ] || exit 3
. "$NVM_DIR/nvm.sh"
nvm "use" "v10.4.1"

out="${LOGDIR}/forever.out"
stdout="${LOGDIR}/swsdvolcal.stdout"
stderr="${LOGDIR}/swsdvolcal.stderr"

logzip() {
  for i in "${out}" "${stdout}" "${stderr}"
  do
    if [ -f "${i}" ]
    then
      n=1
      while [ -f "${i}.${n}.gz" ]
      do
        let n="${n}+1"
      done
      gzip -qc "${i}" > "${i}.${n}.gz"
      cp /dev/null "${i}"
    fi
  done
}

APP="./bin/www"
[ -f "${APP}" ] || exit 3

case "$1" in
start)
  logzip
  forever start   --plain --minUptime=5s --spinSleepTime=5s -a -l "${out}" -o "${stdout}" -e "${stderr}" "${APP}" ;;
stop)
  forever stop    --plain "${APP}" ;;
restart)
  forever restart --plain "${APP}" ;;
*)
  echo "Usage: $SCRIPTNAME {start|stop|restart}" >&2
  exit 4 ;;
esac

:
